# config/stock-analyzer.service
[Unit]
Description=Stock Analysis Assistant
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/stock-analyzer
Environment="PATH=/home/ubuntu/.npm-global/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/home/ubuntu/stock-analyzer/.env
ExecStart=/home/ubuntu/stock-venv/bin/python -m uvicorn server.main:app --host 127.0.0.1 --port 8088
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
